# misc/crontab
59 16 * * *  /home/nccsoft/komutracker/venv/bin/python /home/nccsoft/komutracker/aw-server/report_spent_time.py
30 5 * * *   /home/nccsoft/komutracker/venv/bin/python /home/nccsoft/komutracker/aw-server/reminder_turnon_tracker.py